# scripts/systemd/noahbot-web.service
[Unit]
Description=noahbot playlist-builder web (dev)

[Service]
Type=simple
WorkingDirectory=/home/noah/noahbot
ExecStart=/bin/bash -lc 'exec npm run dev:playlist-web'
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
